# yunwei-sso.service
[Unit]
Description=Yunwei SSO Portal Service - 运维统一门户系统
Documentation=https://github.com/your-repo/yunwei-sso
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=your-username
Group=your-username
WorkingDirectory=/path/to/yunwei_sso/backend
ExecStart=/usr/bin/node src/app.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yunwei-sso

# 环境变量
Environment=NODE_ENV=production
Environment=PORT=3000

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/yunwei_sso/backend/data

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
